# cron/generic/generic_backup_db.cron

## Build backup database, <BACKUP_DATABASE>, hourly
0 6-18 * * * export LD_ASSUME_KERNEL=2.4.1 && <Home_dir>/versions/<version_name>/bin/restore_DB.pl -rebuild -host <BACKUP_HOST> -dbase <BACKUP_DATABASE> -from <SQL_HOST>:<DATABASE> -force -user super_cron -password repus 1> <Data_home_dir>/private/logs/daily_restore_<BACKUP_HOST>.<BACKUP_DATABASE>.log 2> <Data_home_dir>/private/logs/daily_restore_<BACKUP_HOST>.<BACKUP_DATABASE>.err
